# etc/supervisord.conf
[program:botflow]
command=.venv/bin/gunicorn -c etc/gunicorn.conf.py main:app
autostart=true
autorestart=true
numprocs=1
priority = 900
redirect_stderr=true

[supervisorctl]